# systemd/agent-lens-collect.service
[Unit]
Description=Agent Lens — collect Claude Code session traces, then ingest
Documentation=file://__REPO_ROOT__/README.md

[Service]
Type=oneshot
# This host activates Node via mise (a shell hook); systemd user services don't get that hook,
# so node would be off PATH. setup-systemd.sh bakes the resolved node bin dir in here, and runs
# from the repo so any directory-pinned toolchain (mise.toml) resolves correctly.
WorkingDirectory=__REPO_ROOT__
Environment=PATH=__NODE_BIN__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# oneshot ExecStart lines run sequentially; ingest (Stage 2) only runs if collect (Stage 1) succeeds.
ExecStart=__COLLECT_SH__
ExecStart=__INGEST_SH__
# Be a good background citizen.
Nice=10
IOSchedulingClass=idle
# Local-only: no network needed for collection.
