# systemd/agent-lens-collect.timer
[Unit]
Description=Agent Lens — periodic Claude Code trace collection

[Timer]
# A few times a day. Persistent catches runs missed while the machine was off.
OnCalendar=*-*-* 09,13,17,21:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
